# actividad 4/myscript.service
[Unit]
Description=My Script

[Service]
Type=simple
ExecStart=/ruta/a/myscript.sh

[Install]
WantedBy=multi-user.target
